# scene_raytracing/mini-botnet/server_rendering/server.service
[Unit]
Description=client_rendering

[Service]
ExecStart=/root/server_rendering/run
Restart=on-failure

[Install]
WantedBy=multi-user.target

