# scene_raytracing/mini-botnet/client_rendering/client.service
[Unit]
Description=client_rendering

[Service]
ExecStart=/root/client_rendering/run
Restart=on-failure

[Install]
WantedBy=multi-user.target

